# dicom-cleanup.service
[Unit]
Description=DICOM Files Cleanup Service
Documentation=Elimina archivos US (3h), BD (24h) y otros DICOM (30d)

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/DICOMReceiver
ExecStart=/home/ubuntu/DICOMReceiver/cleanup_dicom_files.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=dicom-cleanup

[Install]
WantedBy=multi-user.target
